# assistant-bot.service
[Unit]
Description=Telegram Assistant Bot
After=network-online.target ssh-proxy-tunnel.service
Wants=network-online.target ssh-proxy-tunnel.service

[Service]
Type=simple
WorkingDirectory=/home/andrey/Projects/my-assistant
EnvironmentFile=/home/andrey/Projects/my-assistant/assistant-bot.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/andrey/Projects/my-assistant/bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
